# scripts/systemd/bluesky-bookmark-sync.service
[Unit]
Description=Reader — sync Bluesky bookmarks into the Found feed
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Adjust the path if the repo lives elsewhere on the host.
ExecStart=/usr/bin/node /home/petter/github/reader/scripts/bluesky-bookmark-sync.mjs
# Secrets/state are read from ~/.config/bluesky/ and ~/.config/reader/ —
# no EnvironmentFile needed.
Nice=10
