# scripts/systemd/reader-sources-sync.timer
[Unit]
Description=Reader — linked-sources sync twice daily

[Timer]
# Same slots the retired Sleeper-side X collector used: morning and evening.
# Reddit and HN are free; X bills per post returned, so cadence — not
# volume — is the cost lever.
OnCalendar=*-*-* 07,19:30:00
RandomizedDelaySec=120
Persistent=true

[Install]
WantedBy=timers.target
