# scripts/systemd/reader-hn-candidates.timer
[Unit]
Description=Reader — weekly HN front-page Discover mining

[Timer]
# Weekly is plenty: the collector looks back 30 days, the ingest seam is
# idempotent, and Discover is a slow surface by design.
OnCalendar=Mon 07:10
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
