# scripts/systemd/reader-sources-sync.service
[Unit]
Description=Reader — sync linked sources (X, Reddit, HN) into Found (Worker-side)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Adjust the path if the repo lives elsewhere on the host.
ExecStart=/usr/bin/node /home/petter/github/reader/scripts/sync-sources.mjs
# Auth is read from ~/.config/reader/env (READER_CRON_KEY) —
# no EnvironmentFile needed.
Nice=10
