# scripts/systemd/reddit-saved-sync.service
[Unit]
Description=Reader — sync Reddit saved items into the Found feed
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Adjust the path if the repo lives elsewhere on the host.
ExecStart=/usr/bin/node /home/petter/github/reader/scripts/reddit-saved-sync.mjs
# Secrets/state are read from ~/.config/reddit/ and ~/.config/reader/ —
# no EnvironmentFile needed.
Nice=10
